# php74/.supervisor/supervisord.d/php-worker.ini
[program:php-worker]
process_name=%(program_name)s
command=/opt/remi/php74/root/usr/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
redirect_stderr=true
priority=2
stderr_logfile=/var/log/php.err.log
stdout_logfile=/var/log/php.out.log